# example_cron_job.cron
# NOTE - DO NOT USE THIS IN PRODUCTION, CONFIGURE IT FIRST
* 12 * * * /path/to/backup.py
